# init/psd.service.in
[Unit]
Description=Profile-sync-daemon
Documentation=man:psd(1) man:profile-sync-daemon(1)
Documentation=https://wiki.archlinux.org/index.php/Profile-sync-daemon
Wants=psd-resync.timer
RequiresMountsFor=/home/
After=winbindd.service

[Service]
Slice=background.slice
Type=oneshot
RemainAfterExit=yes
ExecStart=@BINDIR@/psd startup
ExecStartPost=@BINDIR@/psd resync
ExecStop=@BINDIR@/psd unsync
Environment=LAUNCHED_BY_SYSTEMD=1

[Install]
WantedBy=default.target
